# study-mode.service
[Unit]
Description=Block distracting applications while studying
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/study-mode
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target